# config/includes.chroot_after_packages/etc/rc.d/rc.local
#!/bin/sh
echo 'Hangzhou is a city.' >> /hz.txt
